# scripts/neomind.service
[Unit]
Description=NeoMind Server
Documentation=https://github.com/camthink-ai/NeoMind
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=neomind
Group=neomind

# Working directory
WorkingDirectory=/var/lib/neomind

# Executables
ExecStart=/usr/local/bin/neomind-api
ExecReload=/bin/kill -HUP $MAINPID

# Auto-restart
Restart=always
RestartSec=5
StartLimitBurst=10
StartLimitInterval=60

# Environment
Environment="RUST_LOG=info"
Environment="NEOMIND_DATA_DIR=/var/lib/neomind"
Environment="NEOMIND_BIND_ADDR=0.0.0.0:9375"

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/neomind
RestrictRealtime=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX

# Resource limits
LimitNOFILE=65536
MemoryMax=2G

[Install]
WantedBy=multi-user.target
